# harbour-defender.path
[Unit]
Description=Updates Defender hosts file whenever requested
Documentation= man:systemd.path
After=sailfish-unlock-agent
After=pre-user-session.target

[Path]
PathExists=/home/nemo/.config/harbour-defender/update
PathExists=/home/defaultuser/.config/harbour-defender/update

[Install]
WantedBy=user-session.target 
